# agents/burn-agent/systemd/florentine-burn-agent.service.example
[Unit]
Description=Florentine Groteschi Privy burn agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/florentine-groteschi
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run agent:daemon
Restart=always
RestartSec=15
User=florentine-agent
Group=florentine-agent

# Basic hardening. Loosen only if your VPS path layout needs it.
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/florentine-groteschi/agents/burn-agent/logs /opt/florentine-groteschi/agents/burn-agent/state /opt/florentine-groteschi/agents/burn-agent/queue.json

[Install]
WantedBy=multi-user.target
